# invest-agent.service
[Unit]
Description=Invest Agent Scheduler Service
After=network.target

[Service]
# 【关键】设置你的工作目录
WorkingDirectory=/home/ubuntu/invest

# 【关键】指向你的 python 解释器路径和脚本路径
# 如果你用 uv/venv，路径通常是 /home/ubuntu/invest/.venv/bin/python
ExecStart=/home/ubuntu/invest/.venv/bin/python /home/ubuntu/invest/scheduler.py

# 总是自动重启
Restart=always
RestartSec=10

# 设置日志输出不缓冲，方便实时查看
Environment=PYTHONUNBUFFERED=1

# 指定以哪个用户运行（不要用 root，用你的服务器用户名，如 ubuntu）
User=ubuntu

[Install]
WantedBy=multi-user.target
